# autostart/fedora-setup.desktop
[Desktop Entry]
Name=Start Fedora setup script
Type=Application
Exec=HOME/Fedora-Setup-Script/setup.sh
Terminal=true
