# autostart/fedora-setup.desktop
[Desktop Entry]
Name=Start Fedora setup script
Type=Application
Exec=HOME/setup.sh
Terminal=true
